# mailsync.service
[Unit]
Description=Syncing Mail with Remote

[Service]
#User=bhaskar
#Group=bhaskar
ExecStart=/home/bhaskar/bin/mail_sync

